# .root-fs/etc/s6-overlay/s6-rc.d/10-laravel-cache/up
with-contenv
s6-envuidgid docker
if { /opt/php/bin/php /var/www/html/artisan config:cache }
if { /opt/php/bin/php /var/www/html/artisan event:cache }
if { /opt/php/bin/php /var/www/html/artisan route:cache }
/opt/php/bin/php /var/www/html/artisan view:cache
